# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/mysqld_safe

[program:rabbitmq]
command=rabbitmq-server

[program:jetty]
command=/etc/init.d/cloudstack-simulator fstart
